# files/openvpn-crl-refresh.timer
[Unit]
Description=OpenVPN CRL refresh timer

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
